# chef/cookbooks/dobrafaza-singlenode-deploy/templates/default/dobrafaza.conf.erb
description "Dobrafaza upstart script"
author "Paweł Głasek <guasek@gmail.com>"

start on runlevel [2345]
stop on runlevel [016]

respawn

script
    exec /usr/local/bin/grunt --gruntfile /home/dobrafaza/deploy/current/Gruntfile.js serve:dist >> /home/dobrafaza/deploy/current/log 2>&1
end script